# rockpi4/overlays/etc/systemd/system/serial-getty@ttyS2.service.d/override.conf
[Service]
Type=idle
ExecStart=
ExecStart=-/sbin/agetty --autologin debian --noclear %I $TERM
